# packages/rpi-agent/systemd/hanipi-agent.service
[Unit]
Description=HaniPi Sensor Agent
After=network.target

[Service]
Type=simple
User=hanipi
WorkingDirectory=/opt/hanipi
ExecStart=/opt/hanipi/venv/bin/hanipi-agent
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
